# ganache/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[supervisord]
logfile=/var/log/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=trace                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/               ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket

[program:ganache]
command=ganache-cli --blockTime 0.5 --chainId 1440 --networkId 1441 -h 0.0.0.0 -p 8545 -m "ozone access unlock valid olympic save include omit supply green clown session" --db /app/ganache-db --account_keys_path /app/keys.json
stderr_logfile = /dev/stderr
stdout_logfile = /dev/stderr
autorestart = false
stdout_logfile_maxbytes = 0